# supervisor/conf.d/hello.ini
[program:hello]
command=python /home/chang/study/all-kinds-of-examples/supervisor/hello.py &
stdout_logfile=/home/chang/study/all-kinds-of-examples/supervisor/logs/hello.log
stderr_logfile=/home/chang/study/all-kinds-of-examples/supervisor/logs/hello_error.log
stopasgroup=true
killasgroup=true
